# indexer.ini
[program:indexer]
directory=/var/index-service
command=/var/index-service/venv/bin/python run.py
autostart=true
autorestart=true
startretries=3
