# indexer.ini
[program:indexer]
directory=/var/index-service
command=/var/index-service/venv/bin/python index_service.py
autostart=true
autorestart=true
startretries=3
